# etc/supervisord.d/dash.ini
[program:dashboard]
command = /opt/qnib/grafana/bin/dashboard.py -L INFO -l --server
redirect_stderr=true
stdout_logfile=syslog
