# etc/supervisord.d/slurmdash.ini
[program:slurmdash]
command = /opt/qnib/grafana/bin/slurm_dashboard.py -L INFO -l --server
redirect_stderr=true
stdout_logfile=syslog
